# contrib/systemd/kronos-server.service
[Unit]
Description=Kronos control plane
Documentation=https://github.com/ersinkoc/Kronos
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=kronos
Group=kronos
EnvironmentFile=-/etc/kronos/kronos-server.env
ExecStart=/usr/local/bin/kronos server --config /etc/kronos/kronos.yaml
Restart=on-failure
RestartSec=5s
StateDirectory=kronos
RuntimeDirectory=kronos
WorkingDirectory=/var/lib/kronos
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/kronos
ReadOnlyPaths=/etc/kronos
CapabilityBoundingSet=
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
